# specs/ermete-system-services/SOURCES/usr/lib/systemd/user/ermete-wallpaper.service
[Unit]
Description=Ermete OS Wallpaper Service
Documentation=https://github.com/patapem/ermete-forge
BindsTo=graphical-session.target
After=graphical-session-pre.target

[Service]
Type=simple
Slice=session.slice
ExecStart=/usr/bin/bash -c 'if [ -x "$(command -v swaybg)" ]; then if [ -f "$HOME/.config/ermete/wallpaper.png" ]; then exec swaybg -i "$HOME/.config/ermete/wallpaper.png" -m fill; elif [ -f "/usr/share/backgrounds/ermete/default.png" ]; then exec swaybg -i "/usr/share/backgrounds/ermete/default.png" -m fill; else exec swaybg -c "#1e1e2e"; fi; else sleep infinity; fi'
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=niri-session.target graphical-session.target
